# docker/kafka/configs/topic.conf
[program:topicoPag]
command=/usr/bin/criar_topico.sh